# templates/init.conf
#!upstart
 
description "{{ app_name }}"
author      "{{ user_name }}"
 
start on (local-filesystems and net-device-up IFACE=eth0) 
stop  on shutdown
 
respawn                # restart when job dies
respawn limit 5 60     # give up restart after 5 respawns in 60 seconds

env NODE_ENV=production

export NODE_ENV

script
  exec node /var/www/{{ app_name }}/live/{{ app_name }}/app.js >> /var/log/{{ app_name }}.log 2>&1
end script